# crontab
# /etc/cron.d/cactibackup
0 1 * * * root /usr/local/bin/runCactiBackup.sh